# templates/etc/systemd/td-agent.service.erb
[Unit]
Description=<%= package_name %>: Fluentd based data collector for Treasure Data
Documentation=https://docs.treasuredata.com/articles/td-agent
After=network-online.target
Wants=network-online.target

[Service]
User=<%= Shellwords.shellescape(package_name) %>
Group=<%= Shellwords.shellescape(package_name) %>
LimitNOFILE=65536
# XXX: What version should we use jemalloc?
# Environment=LD_PRELOAD=<%= package_dir_opt %>/embedded/lib/libjemalloc.so
Environment=GEM_HOME=<%= package_dir_opt %>/
Environment=GEM_PATH=<%= package_dir_opt %>/
Environment=FLUENT_CONF=/etc/<%= package_name %>/<%= package_name %>.conf
Environment=FLUENT_PLUGIN=/etc/<%= package_name %>/plugin
Environment=FLUENT_SOCKET=/var/run/<%= package_name %>/<%= package_name %>.sock
Environment=TD_AGENT_OPTIONS=
<% if pkg_type == 'deb' %>
EnvironmentFile=-/etc/default/<%= package_name %>
<% else %>
EnvironmentFile=-/etc/sysconfig/<%= package_name %>
<% end %>
PIDFile=<%= Shellwords.shellescape(File.join(root_dir, "var", "run", package_name, "#{package_name}.pid")) %>
RuntimeDirectory=<%= Shellwords.shellescape(package_name) %>
Type=forking
ExecStart=/opt/td-agent/bin/fluentd --log <%= Shellwords.shellescape(File.join(root_dir, "var", "log", package_name, "#{package_name}.log")) %> --daemon <%= Shellwords.shellescape(File.join(root_dir, "var", "run", package_name, "#{package_name}.pid")) %> $TD_AGENT_OPTIONS
ExecStop=/bin/kill -TERM ${MAINPID}
ExecReload=/bin/kill -HUP ${MAINPID}
Restart=always
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
